# ai-core/hybrid-engine.bak.20251112231325/systemd/cockswain-bridge.service
[Unit]
Description=Cockswain Bridge Service
After=network.target
Requires=network.target

[Service]
Type=simple
User=cocksmain
WorkingDirectory=/srv/cockswain-core/ai-core/hybrid-engine
ExecStart=/usr/bin/python3 -m hybrid_engine.bridge.bridge_service
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
PYTHONPATH=/srv/cockswain-core/ai-core/hybrid-engine

[Install]
WantedBy=multi-user.target
