# ai-core/mind-proxy-v0.1/deploy/systemd/cockswain-mind-proxy.service
[Unit]
Description=Cockswain Mind Proxy Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cocksmain
WorkingDirectory=/srv/cockswain-core/ai-core/mind-proxy
EnvironmentFile=/srv/cockswain-core/ai-core/mind-proxy/.env
ExecStart=/srv/cockswain-core/ai-core/mind-proxy/.venv/bin/mind-proxy
Restart=always
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
